# rootfs/etc/s6/nginx/finish
#!/bin/bash
echo "> stopping nginx"
s6-svscanctl -t /etc/s6
